# src/configs/main/systemd/user/cleanup-old-files.timer
[Unit]
Description=Daily cleanup of old files

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
